# conf/etc/cron.d/voron-cron
# m h  dom mon dow   user	command

# Synchro des fichiers
@reboot pi bash /home/pi/scripts/synchro-file.sh

# Initialisation des GPIO
@reboot pi bash /home/pi/scripts/init-gpio.sh

# Btn RPI
@reboot root python /home/pi/scripts/btn-rpi.py

# Upgrade des scripts (en root)
@reboot root bash /home/pi/scripts/update-script-root.sh
